# conf/conf.d/novnc.sv.conf
[program:novnc]
command=/usr/bin/noVNC/utils/launch.sh --vnc localhost:5900
autorestart=true
user=user
priority=100
