# conf/conf.d/openbox.sv.conf
[program:openbox]
environment=HOME="/home/user",DISPLAY=":1",USER="user"
command=/usr/bin/openbox
user=user
autorestart=true
priority=300
